# package/root/etc/update-motd.d/11-rock64-help-text
#!/bin/sh

if [ $(id -u) -eq 0 ]; then
    echo ""
    echo "Update system with:"
    echo ""
    echo " * rock64_update_package.sh"
fi
