# package/root/etc/cron.d/make_nas_processes_faster
* * * * * root for i in `pgrep "ftpd|nfsiod|smbd|afpd|cnid"` ; do ionice -c1 -p $i ${XU4_HMP_Fix}; done &>/dev/null
