# package/root/etc/cron.d/startup
@reboot root /usr/local/sbin/rock64_fix_performance.sh

